# .local/scripts/etc_systemd_system/opt.mount
[Unit]
Description=Bind /home/opt to /opt

[Mount]
What=/home/opt
Where=opt
Type=none
Options=bind

[Install]
WantedBy=multi-user.target

