# .local/scripts/etc_systemd_system/usr.mount
[Unit]
Description=Bind /home/usr to /usr

[Mount]
What=/home/usr
Where=usr
Type=none
Options=bind

[Install]
WantedBy=multi-user.target

